# contrib/osec.timer
[Unit]
Description=Run OSEC every day at midnight

[Timer]
OnCalendar=*-*-* 3:00
RandomizedDelaySec=60m

[Install]
WantedBy=multi-user.target
