# contrib/osec.service
[Unit]
Description=OSEC integrity check

[Service]
Type=oneshot
ExecStart=/usr/share/osec/osec.cron
